# botrss-streamlit.service
[Unit]
Description=BotRSS Streamlit UI
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/rssgenerator
ExecStart=/rssgenerator/.venv/bin/streamlit run app.py --server.port 8501 --server.headless true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 