# botrss-api.service
[Unit]
Description=BotRSS FastAPI Server
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/root/rssgenerator
ExecStart=/root/rssgenerator/.venv/bin/uvicorn rss_api:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 